# deploy/rt-handler.service
[Unit]
Description=RT Handler - real-time GPIO monitor
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/rt-handler
ExecStart=/usr/sbin/rt-handler -b ${BOARD}
Restart=on-failure
RestartSec=5
LimitRTPRIO=99
ProtectSystem=strict
ReadWritePaths=/dev /sys /proc /etc/sysconfig
DeviceAllow=char-gpiochip rw
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
